# debian/cron.weekly
#! /bin/sh

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
MAINTENANCE=/usr/bin/http-replicator_maintenance
NAME=http-replicator

if [ -x $MAINTENANCE ] && [ -f /etc/default/$NAME ] ; then
	source /etc/default/$NAME
else
	exit 0
fi

set -e

$MAINTENANCE $MAINTENANCE_OPTS

exit 0
